# service/dbus-systemcalc-py/run
#!/bin/sh
exec 2>&1
exec /opt/victronenergy/dbus-systemcalc-py/dbus_systemcalc.py
